# infra/systemd/ploy-claimer.service
[Unit]
Description=Ploy Auto Claimer (Polymarket Redeem)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=ploy
Group=ploy
WorkingDirectory=/opt/ploy

# Keep logs lean by default; debug-level logs can exceed disk limits in hours.
Environment=RUST_LOG=info,ploy=info,sqlx=warn
Environment=LOG_DIR=/opt/ploy/logs
EnvironmentFile=-/opt/ploy/.env
EnvironmentFile=-/opt/ploy/env/claimer.env

# Redeem resolved winning positions periodically.
# NOTE: requires POLYMARKET_PRIVATE_KEY/PRIVATE_KEY in env files above.
ExecStart=/opt/ploy/bin/ploy claim --interval 60 --min-size 1

Restart=always
RestartSec=15

NoNewPrivileges=true
ProtectSystem=full
ProtectHome=true
ReadWritePaths=/opt/ploy/data /opt/ploy/logs /opt/ploy/config /opt/ploy/env
PrivateTmp=true

MemoryMax=512M
CPUQuota=50%
LimitNOFILE=65535

StandardOutput=journal
StandardError=journal
SyslogIdentifier=ploy-claimer

[Install]
WantedBy=multi-user.target

